# orchestration/launchd/com.ghostbrain.profile-monthly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Monthly profile decay + Stable promotion. Runs the 1st of every month
  at 22:00 local time per SPEC §10.2.

  This file is a TEMPLATE. Replace placeholders before loading.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ghostbrain.profile-monthly</string>

    <key>ProgramArguments</key>
    <array>
        <string>__REPO_ROOT__/.venv/bin/python</string>
        <string>-m</string>
        <string>ghostbrain.profile.decay</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_ROOT__</string>

    <key>RunAtLoad</key>
    <false/>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__REPO_ROOT__/logs/profile-monthly.log</string>

    <key>StandardErrorPath</key>
    <string>__REPO_ROOT__/logs/profile-monthly.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>VAULT_PATH</key>
        <string>__VAULT_PATH__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
